# deploy/systemd/steam-notifier.service
# Steam sale digest — one-shot. Triggered by steam-notifier.timer.
# Install: see deploy/systemd/README.md. Edit the path/User below if your
# checkout or username differs, then `sudo systemctl daemon-reload`.

[Unit]
Description=Steam sale daily digest -> Slack
Documentation=https://github.com/your/repo
Wants=network-online.target
After=network-online.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=vuksan
WorkingDirectory=/home/vuksan/steam-sale-notifier
# `run --rm` starts a fresh container, runs once, and is removed on exit.
ExecStart=/usr/bin/docker compose -f /home/vuksan/steam-sale-notifier/docker-compose.yml run --rm notifier
